# dist/systemd/cosmic-toysd.service
[Unit]
Description=COSMIC color-picker daemon
Documentation=https://github.com/Pyxyll/cosmic-toys
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/cosmic-toysd
Restart=on-failure
RestartSec=5
# Crash logs land in journalctl --user -u cosmic-toysd
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
